# supervisord.conf
[supervisord]
nodaemon=true

[program:celery]
command=celery worker -A gitsniffer.tasks --loglevel=DEBUG
user=nobody
startsecs=10
stopwaitsecs=600
numprocs=1
stdout_logfile=/var/log/celery/worker.log
stderr_logfile=/var/log/celery/worker.log

[program:rethinkdb]
command=rethinkdb --bind all
startsecs=60

[program:redis]
command=/usr/sbin/redis-server

[program:gitsniffer]
command=/usr/bin/gitsniffer_target
stdout_logfile=/var/log/gitsniffer.log
stderr_logfile=/var/log/gitsniffer.log
